# runner.cron
*/5 * * * * user flock -w 3600 /var/tmp/runner -c /home/user/runner-wrapper.sh